# etc/cron.d/backup
* * * * * hoyoswelch /home/hoyoswelch/hoyoswelch/backup.sh
