# etc/cron.d/backup
1 * * * * $HOME/hoyoswelch/hoyoswelch/backup
