# etc/cron.d/backup
00 00 * * * hoyoswelch /home/hoyoswelch/hoyoswelch/backup.sh
